# nominatim/nominatim.socket
[Unit]
Description=Gunicorn socket for Nominatim

[Socket]
ListenStream=/run/nominatim.sock
SocketUser=http

[Install]
WantedBy=multi-user.target
