# lcdproc/lcdd.service
[Unit]
Description=LCD display driver

[Service]
Type=forking
ExecStart=/usr/bin/LCDd -c /etc/LCDd.conf

[Install]
WantedBy=multi-user.target
