# openrgb/openrgb.service
[Unit]
Description=Run openrgb server
After=network.target lm_sensors.service

[Service]
ExecStart=/usr/bin/openrgb --server
Restart=on-failure

[Install]
WantedBy=multi-user.target
