# openrgb/openrgb.service
[Unit]
Description=Run OpenRGB server
After=network.target lm_sensors.service

[Service]
ExecStart=/usr/bin/openrgb --server --config /etc/openrgb
Restart=on-failure

[Install]
WantedBy=multi-user.target
